# .deploy/forgedesk-worker.conf
# ForgeDesk Studio - Queue Worker Supervisor Configuration
# Place in /etc/supervisor/conf.d/forgedesk-worker.conf
# Then run: sudo supervisorctl reread && sudo supervisorctl update

[program:forgedesk-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/forgedesk/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --rest=1
directory=/var/www/forgedesk
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/forgedesk-worker.log
stopwaitsecs=3600
